# cronus/upstart.conf
description     "Cronus Agent"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

# First define what is ROOT directory and USER to run as

script

	if [ -f /usr/bin/setuidgid ]; then
		cd PATH_TO_PACKAGE
		exec setuidgid APP_USER cronus/scripts/startup nodaemon
	else
		exec  su - APP_USER -c "PATH_TO_PACKAGE/cronus/scripts/startup nodaemon"
	fi

end script

post-stop exec su - APP_USER -c PATH_TO_PACKAGE/cronus/scripts/shutdown



